# packaging/debian/postinst
#!/usr/bin/sh

set -e

case "$1" in
    configure)
        # Add shebang to the daemon script
        sed -i '1s|^|#!/opt/better-backlight/venv/bin/python\n|' /opt/better-backlight/better_backlight/entrypoints/daemon.py

        export PIPENV_VERBOSITY=-1
        cd /opt/better-backlight/

        python3 -m venv /opt/better-backlight/venv

        mkdir -p /opt/better-backlight/bin
        ln -s /opt/better-backlight/better_backlight/entrypoints/daemon.py /opt/better-backlight/bin/better-backlight
        chmod +x /opt/better-backlight/bin/better-backlight

        # Dependencies installation. Idempotent.
        . /opt/better-backlight/venv/bin/activate

        pip install pipenv
        pipenv install

        # Enable and start the service
        systemctl daemon-reload
        systemctl enable --now better-backlight.service
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
